# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:colladaconv]
directory=/app
command=python /app/server.py

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
